# deployment/power-market-collector.service
[Unit]
Description=Power Market Real-time Data Collector
After=network.target postgresql.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/power_market_pipeline
Environment="PATH=/home/ubuntu/.local/bin:/usr/local/bin:/usr/bin:/bin"
Environment="DATABASE_URL=postgresql://localhost:5432/power_market"
ExecStart=/home/ubuntu/.local/bin/python /home/ubuntu/power_market_pipeline/realtime_collector.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=power-market-collector

# Performance and security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/ubuntu/power_market_pipeline/data

[Install]
WantedBy=multi-user.target